# deploy/dolOS.service
[Unit]
Description=dolOS AI Agent
Documentation=https://github.com/your-repo/dolOS
After=network-online.target
Wants=network-online.target
# Uncomment if Ollama runs as a separate systemd service:
# After=ollama.service
# Prevent runaway restart loops: max 5 restarts within 2 minutes
StartLimitIntervalSec=120s
StartLimitBurst=5

[Service]
Type=simple
# By default runs as root. To run as a dedicated user, create it first:
#   sudo useradd -r -s /bin/false -d /opt/dolOS dolos
# then uncomment the next two lines:
# User=dolos
# Group=dolos
WorkingDirectory=/opt/dolOS
EnvironmentFile=/opt/dolOS/.env

# Build the React dashboard before starting (skips if ui/dist already exists and is current)
ExecStartPre=/bin/bash -c 'cd /opt/dolOS/ui && npm install --silent && npm run build'
# Install Playwright Chromium browser if not already present (idempotent)
# --yes is an npx flag (auto-confirms package install); "playwright install chromium" is the command.
# This installs the Node.js Playwright Chromium browser used by @playwright/mcp.
# Idempotent: skips if already installed. || true prevents blocking if Node.js is absent.
ExecStartPre=/bin/bash -c 'npx --yes playwright install chromium 2>/dev/null || true'

ExecStart=/opt/dolOS/.venv/bin/python main.py

# Restart on crash; don't restart on clean exit (exit code 0)
Restart=on-failure
RestartSec=10s

# Give the agent 15 seconds to finish in-flight requests on stop
TimeoutStopSec=15s
KillMode=mixed

# Security hardening (re-enable when not running as root)
# NoNewPrivileges=true
# PrivateTmp=true

# journald captures stdout/stderr automatically — no log files needed
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dolOS

[Install]
WantedBy=multi-user.target
